# cloudify_marketplace/resources/cloudify_security_reconfiguration.cron
* * * * * root /opt/cfy_cli/bin/python /opt/cloudify_security_reconfiguration
